# rootfs/usr/systemd/system/sddm.service.d/boot-order.conf
[Unit]
# Ensure SDDM starts only after user sessions and dbus are ready
After=systemd-user-sessions.service dbus.service
Requires=dbus.service

