# rootfs/usr/lib/systemd/system/nix-daemon.service.d/override.conf
[Unit]
# Ensure all tmpfiles (including the /home bind mount) are set up
# before nix-daemon starts, so that /home is a real directory.
After=systemd-tmpfiles-setup.service
Requires=systemd-tmpfiles-setup.service
